# deploy/systemd/resume-ai-auto-update.timer
[Unit]
Description=Run Resume AI auto update every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s
Persistent=true
Unit=resume-ai-auto-update.service

[Install]
WantedBy=timers.target
